# autostart/Dropbox.desktop
[Desktop Entry]
Type=Application
Name=Dropbox
Exec=/home/jpreiss/.dropbox-dist/dropboxd
