# autostart/Hexchat.desktop
[Desktop Entry]
Type=Application
Name=Hexchat
Exec=hexchat

