# autostart/Pidgin.desktop
[Desktop Entry]
Type=Application
Name=Pidgin
Exec=pidgin