# autostart/Skype.desktop
[Desktop Entry]
Type=Application
Name=Skype
Exec=skype