# autostart/Twitter.desktop
[Desktop Entry]
Type=Application
Name=Twitter
Exec=corebird
