# etc/unix/autostart/devilspie.desktop
[Desktop Entry]
Exec=devilspie
Name=Devilspie
Type=Application
X-GNOME-Autostart-enabled=true
