# etc/unix/autostart/Hexchat.desktop
[Desktop Entry]
Exec=hexchat
Name=HexChat
Type=Application
Hidden=false

