# etc/unix/autostart/Launchy.desktop
[Desktop Entry]
Type=Application
Name=Launchy
Exec=launchy

