# etc/unix/autostart/Pidgin.desktop
[Desktop Entry]
Type=Application
Name=Pidgin
Exec=pidgin
Terminal=false
X-GNOME-Autostart-enabled=true

