# etc/unix/autostart/Skype.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]                                                                 
Type=Application
Name=Skype
Exec=skype
