# etc/unix/autostart/Synergy.desktop
[Desktop Entry]
Type=Application
Name=Synergy
Exec=synergyc 192.168.1.10
