# etc/unix/autostart/Thunderbird.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Exec=thunderbird
Name=Thunderbird
Type=Application
