# etc/unix/autostart/Thunderbird.desktop
[Desktop Entry]
Type=Application
Name=Thunderbird
Exec=thunderbird

