# etc/unix/autostart/Thunderbird.desktop
[Desktop Entry]
Exec=thunderbird
Name=Thunderbird
Type=Application
Hidden=false

