# etc/unix/autostart/Twitter.desktop
[Desktop Entry]
Exec=cowbird
Name=Twitter
Type=Application
