# etc/unix/autostart/Twitter.desktop
[Desktop Entry]
Exec=corebird
Name=Twitter
Type=Application
